# Library/LaunchAgents/com.flaport.jupyterhub.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.flaport.jupyterhub</string>

  <key>EnvironmentVariables</key>
  <dict>
      <key>PATH</key>
      <string>/Users/flaport/.local/share/uv/tools/jupyterhub/bin:/Users/flaport/.local/bin:/opt/homebrew/bin:/usr/local/bin:/System/Cryptexes/App/usr/bin:/usr/bin:/bin:/usr/sbin:/sbin:/var/run/com.apple.security.cryptexd/codex.system/bootstrap/usr/local/bin:/var/run/com.apple.security.cryptexd/codex.system/bootstrap/usr/bin:/var/run/com.apple.security.cryptexd/codex.system/bootstrap/usr/appleinternal/bin:/Users/flaport/.cargo/bin:/Applications/Ghostty.app/Contents/MacOS</string>
      <key>HOME</key>
      <string>/Users/flaport</string>
  </dict>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/flaport/.local/share/uv/tools/jupyterhub/bin/jupyter</string>
    <string>lab</string>
    <string>--port</string>
    <string>8000</string>
    <string>--no-browser</string>
  </array>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/flaport/.jupyter/jupyterhub.out.log</string>

  <key>WorkingDirectory</key>
  <string>/Users/flaport</string>

  <key>StandardErrorPath</key>
  <string>/Users/flaport/.jupyter/jupyterhub.err.log</string>
</dict>
</plist>
